# src/main/dist/etc/upsilon-drone.service
[Unit]
Description=upsilon-drone
After=syslog.target network.target remote-fs.target nss-lookup.target

[Service]
User=upsilon
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=upsilon-drone
ExecStart=/usr/share/upsilon-drone/bin/upsilon-drone

[Install]  
WantedBy=multi-user.target
